# backend/config/mes.service
[Unit]
Description=MES Spring Boot Application
Documentation=https://github.com/your-repo/smart-factory-mes
After=syslog.target network.target mysql.service
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/mes/smart-factory-mes/backend

# Java 실행 명령
ExecStart=/usr/bin/java \
    -Xms512m \
    -Xmx1024m \
    -XX:+UseG1GC \
    -XX:MaxGCPauseMillis=200 \
    -XX:+HeapDumpOnOutOfMemoryError \
    -XX:HeapDumpPath=/var/log/mes/heap-dump.hprof \
    -Djava.security.egd=file:/dev/./urandom \
    -Dspring.profiles.active=prod \
    -Dserver.port=8080 \
    -jar target/mes-inno-0.0.1-SNAPSHOT.jar

# 프로세스 종료 시그널
ExecStop=/bin/kill -TERM $MAINPID
SuccessExitStatus=143
TimeoutStopSec=10

# 재시작 정책
Restart=on-failure
RestartSec=5
StartLimitInterval=600
StartLimitBurst=3

# 환경 변수
EnvironmentFile=-/etc/environment.mes
Environment="SPRING_PROFILES_ACTIVE=prod"
Environment="TZ=Asia/Seoul"

# 리소스 제한
LimitNOFILE=65535
LimitNPROC=4096

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mes

# 보안 강화 옵션
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mes /var/log/mes

[Install]
WantedBy=multi-user.target